# www.zope.org/trunk/etc/cron.d/zeopack
# Pack the zope.org storage every Sunday at 4:43 AM
43 4 * * 0      zope    /var/zope/www.zope.org/bin/zeopack
