# www.zope.org/trunk/etc/supervisord.conf
[unix_http_server]
file=%(here)s/../var/supervisor.sock
chmod=0600

[supervisord]
logfile=%(here)s/../var/log/supervisord.log
logfile_maxbytes=5MB
logfile_backups=10
loglevel=info
pidfile=%(here)s/../var/supervisord.pid ;
childlogdir=%(here)s/../var/log
nodaemon=false              ; (start in foreground if true;default false)
minfds=1024                 ; (min. avail startup file descriptors;default 1024)
minprocs=200                ; (min. avail process descriptors;default 200)
directory=/srv/www.zope.org
user=plone

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix://%(here)s/../var/supervisor.sock

[program:zeo]
command = %(here)s/../bin/zeo fg
autostart= true
autorestart = true

[program:zeoclient]
command = %(here)s/../bin/instance console
priority = 2
redirect_stderr = true
autostart= true
autorestart = true

