# www.zope.org/trunk/etc/supervisord.conf
[unix_http_server]
file=/srv/www.zope.org/var/supervisor.sock
chmod=0600

[supervisord]
logfile=/srv/www.zope.org/var/log/supervisord.log
logfile_maxbytes=5MB
logfile_backups=10
loglevel=info
pidfile=/srv/www.zope.org/var/supervisord.pid ;
childlogdir=/srv/www.zope.org/var/log
nodaemon=false              ; (start in foreground if true;default false)
minfds=1024                 ; (min. avail startup file descriptors;default 1024)
minprocs=200                ; (min. avail process descriptors;default 200)
directory=/srv/www.zope.org
user=plone

[program:zeo]
command = /srv/www.zope.org/bin/zeo fg
autostart= true
autorestart = true

[program:zeoclient]
command = /srv/www.zope.org/bin/instance console
priority = 2
redirect_stderr = true
autostart= true
autorestart = true

[supervisorctl]
serverurl=unix:///srv/www.zope.org/var/supervisor.sock

