# www.zope.org/trunk/etc/supervisor.conf
[unix_http_server]
file=%(here)/var/supervisor.sock
chmod=0600

[supervisord]
logfile=%(here)/var/log/supervisord.log
logfile_maxbytes=5MB
logfile_backups=10
loglevel=info
pidfile=%(here)/var/supervisord.pid ;
childlogdir=%(here)/var
nodaemon=false              ; (start in foreground if true;default false)
minfds=1024                 ; (min. avail startup file descriptors;default 1024)
minprocs=200                ; (min. avail process descriptors;default 200)
directory=%(here)

[program:zeo]
command = %(here)/bin/zeo fg
autostart= true
autorestart = true

[program:zeoclient]
command = %(here)/bin/instance console
priority = 2
redirect_stderr = true
autostart= true
autorestart = true

